# conf/crontab
*/5  *  *  *  * sh /root/srcs/s/shells/detect-stocks.sh
0 3 * * 0 bash /usr/bin/ssl_update.sh
22 0 * * * "/root/.acme.sh"/acme.sh --cron --home "/root/.acme.sh" > /dev/null